# systemd/system/monitoringLLD.target
#
# THIS FILE IS AUTOGENERATED WITH ANSIBLE!
# Do not edit it directly
#
[Unit]
Description=List daemons for dynamic monitoring in %p group
After=network-online.target multi-user.target
