# systemd/system/onfailureLLD@.service
#
# THIS FILE IS AUTOGENERATED WITH ANSIBLE!
# Do not edit it directly
#
[Unit]
Description=Service send alert when %i entered fail state

[Service]
Type=oneshot
ExecStart=-/usr/bin/zabbix_sender -k monitoringLLD.[%i.state] -o 0 -c /etc/zabbix/zabbix_agentd.conf
