# systemd/system/probeLLD@.service
#
# THIS FILE IS AUTOGENERATED WITH ANSIBLE!
# Do not edit it directly
#
[Unit]
Description=Just plug for adding %i to zabbix via LLDMonitoring.target
PartOf=monitoringLLD.target

[Service]
Type=oneshot
ExecStart=/bin/echo %n started
