# systemd/system/td-sip-api.service
#
# THIS FILE IS AUTOGENERATED WITH ANSIBLE!
# Do not edit it directly
#
[Unit]
Description=API for sip-server on televoip
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
User=root
Group=root

Environment="VOIP_BIN=/usr/local/voip/bin"
Environment="SCRIPT_NAME=td-sip-api/run"
Type=simple

Restart=always

ExecStart=/usr/bin/env python3 ${VOIP_BIN}/${SCRIPT_NAME}.py
ExecStop=/bin/echo "Stopping %p"


